# scripts/ticketing.service
[Unit]
Description=Iranmehr Ticketing System
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/imehrTicketing
Environment="PATH=/path/to/imehrTicketing/venv/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/path/to/imehrTicketing/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ticketing

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/path/to/imehrTicketing/storage /path/to/imehrTicketing/logs /var/log/ticketing

# Resource limits
LimitNOFILE=65536
MemoryMax=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target

